# scripts/systemd/substrate.service
[Unit]
Description=Substrate AI Agent Orchestration Service
After=network.target
Documentation=https://github.com/rookdaemon/substrate
OnFailure=substrate-recovery.service
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
User=rook
Group=rook
WorkingDirectory=/home/rook/substrate/server
Environment="NODE_ENV=production"
Environment="HOME=/home/rook"
# PATH: Uses default node symlink for version-agnostic deployment
Environment="PATH=/home/rook/.nvm/versions/node/default/bin:/home/rook/.local/bin:/usr/local/bin:/usr/bin:/bin"

# Start the substrate server (uses supervisor.js for auto-rebuild on exit code 75)
ExecStart=/usr/bin/env npm run start

# Reset recovery attempt counter on successful start
ExecStartPost=/home/rook/substrate/scripts/recovery.sh --reset

# Restart policy
Restart=on-failure
RestartSec=10s

# Resource limits (adjust based on your deployment needs)
# MemoryMax: 2GB suitable for typical substrate workloads
# CPUQuota: 200% allows use of 2 CPU cores (adjust for your hardware)
LimitNOFILE=65536
MemoryMax=2G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=substrate

[Install]
WantedBy=multi-user.target
